# audio-processor/systemd/audio-processor-api.service
[Unit]
Description=Audio Processor API Service
After=network.target postgresql.service redis-server.service

[Service]
Type=notify
User=audioprocessor
Group=audioprocessor
WorkingDirectory=/home/audioprocessor/audio-processor
Environment="PATH=/home/audioprocessor/audio-processor/venv/bin"
ExecStart=/home/audioprocessor/audio-processor/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
